# crontab
*/10 *  * * *   root    /local/scratch/nyancat/reset.sh
